# buildfiles/fedora/lib/systemd/system/appserver-watcher.service

[Unit]
Description=appserver watcher

[Service]
ExecStart=/opt/appserver/bin/appserver-watcher -D FOREGROUND

[Install]
WantedBy=multi-user.target
