# buildfiles/fedora/lib/systemd/system/appserver.service

[Unit]
Description=appserver

[Service]
ExecStart=/opt/appserver/bin/appserver -D FOREGROUND

[Install]
WantedBy=multi-user.target
